# kiosk.timer
[Unit]
Description=Kiosk Service

[Timer]
OnCalendar=*-*-* 07:45:00

[Install]
WantedBy=timers.target
